# ec2/systemd/binance-dca.timer
[Unit]
Description=Daily Binance DCA Order Timer
Requires=binance-dca.service

[Timer]
OnCalendar=daily
# Run at 6:00 AM UTC daily
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target